# systemd/github-analysis-status-publisher.timer
[Unit]
Description=Publish public GitHub analysis status summary twice daily

[Timer]
OnCalendar=*-*-* 08:50:00
OnCalendar=*-*-* 20:50:00
Persistent=true

[Install]
WantedBy=timers.target
