# systemd/github-analysis-status-publisher.timer
[Unit]
Description=Publish public GitHub analysis status summary hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
